# deploy/systemd/religio-send-daily-digest.timer
[Unit]
Description=Religio Pro — daily 07:00 activity digest fan-out
Documentation=https://github.com/galihsidik86/travel/blob/main/deploy/DEPLOYMENT.md

[Timer]
# Run daily at 07:00 local time — owner reads it with morning coffee.
OnCalendar=*-*-* 07:00:00
# If the host was down at trigger time, run on next boot.
Persistent=true
# Add up to 5 min jitter so a fleet of hosts doesn't hammer the DB at once.
RandomizedDelaySec=5min
Unit=religio-send-daily-digest.service

[Install]
WantedBy=timers.target
