# deploy/systemd/religio-send-landing-slow.service
[Unit]
Description=Religio Pro — daily landing-page speed budget check
Documentation=https://github.com/galihsidik86/travel/blob/main/deploy/DEPLOYMENT.md
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=oneshot
User=religio
Group=religio
WorkingDirectory=/srv/religio-pro
EnvironmentFile=/srv/religio-pro/.env
ExecStart=/usr/bin/npm run job:send-landing-slow
StandardOutput=append:/var/log/religio/landing-slow.log
StandardError=append:/var/log/religio/landing-slow.log

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/religio /srv/religio-pro/private
