# deploy/systemd/religio-expire-docs.service
[Unit]
Description=Religio Pro — expire overdue JemaahDocument rows
Documentation=https://github.com/galihsidik86/travel/blob/main/deploy/DEPLOYMENT.md
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=oneshot
User=religio
Group=religio
WorkingDirectory=/srv/religio-pro
EnvironmentFile=/srv/religio-pro/.env
ExecStart=/usr/bin/npm run job:expire-docs
StandardOutput=append:/var/log/religio/expire-docs.log
StandardError=append:/var/log/religio/expire-docs.log

# Sandboxing — least privilege
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/religio /srv/religio-pro/private
