# deploy/systemd/religio-send-task-overdue.service
[Unit]
Description=Religio Pro — daily overdue task escalation
Documentation=https://github.com/galihsidik86/travel/blob/main/deploy/DEPLOYMENT.md
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=oneshot
User=religio
Group=religio
WorkingDirectory=/srv/religio-pro
EnvironmentFile=/srv/religio-pro/.env
ExecStart=/usr/bin/npm run job:send-task-overdue
StandardOutput=append:/var/log/religio/task-overdue.log
StandardError=append:/var/log/religio/task-overdue.log

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/religio
